# app/resonant-laboratory/ansible/xenial/girder.service.template
[Unit]
Description=Girder
After=network.target

[Service]
Environment="PYTHONPATH=${PYTHONPATH}:{{ storage }}/girder"
ExecStart=/usr/bin/python -m girder
